# crontab
* * * * * cd <path_to_LogSignature_installed> && <path_to_python3> <path_to_LogSignature/main.py>
